# build/bgerp/scripts/bgerp.service
[Unit]
Description=BGERP Enterprise Process Management System
After=network.target
 
[Service]
Type=forking
ExecStart=/opt/BGERP/erp.sh start
ExecStop=/opt/BGERP/erp.sh stop
 
[Install]
WantedBy=multi-user.target