# srcx/deploy/demo/crontab
10  */3   *   *   *   /home/demo/demo_reset.sh 2>&1 >/home/demo/demo_reset.log
